# projects/codex-unified-history/launchagents/com.example.codex-unified-history.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.codex-unified-history</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__HOME__/.codex/unified-history/codex_unified_history_agent.sh</string>
  </array>
  <key>StandardOutPath</key>
  <string>__HOME__/.codex/unified-history/launchd.out</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.codex/unified-history/launchd.err</string>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>WatchPaths</key>
  <array>
    <string>__HOME__/.cc-switch/cc-switch.db</string>
    <string>__HOME__/.codex/config.toml</string>
    <string>__HOME__/.codex/state_5.sqlite</string>
  </array>
</dict>
</plist>

